# deploy/supervisor/hedge.conf
; Canonical v7.9 supervisor config
; NOTE: Production uses /etc/supervisor/conf.d/hedge*.conf (app.py).
;       This repo config is kept in sync for reference.

[program:hedge-executor]
directory=/root/hedge-fund
command=/root/hedge-fund/venv/bin/python -m execution.executor_live
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
environment=ENV="prod",BINANCE_TESTNET="1",PYTHONUNBUFFERED="1",ALLOW_OPEN_POSITIONS="1",PREDICTION_PHASE="P2_PRODUCTION"

[program:hedge-sync_state]
directory=/root/hedge-fund
command=/root/hedge-fund/venv/bin/python execution/sync_state.py
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
environment=ENV="prod",BINANCE_TESTNET="1",PYTHONUNBUFFERED="1"

[program:hedge-dashboard]
directory=/root/hedge-fund
command=/root/hedge-fund/venv/bin/python -m streamlit run dashboard/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
environment=ENV="prod",BINANCE_TESTNET="1",PYTHONUNBUFFERED="1"

; Legacy pipeline processes — no longer used in v7.6+
; [program:hedge-pipeline-shadow-heartbeat]
; [program:hedge-pipeline-compare]

[group:hedge]
programs=hedge-executor,hedge-sync_state,hedge-dashboard
