# deploy/supervisor/hedge-dashboard.conf
[program:hedge-dashboard]
directory=/root/hedge-fund
command=/root/hedge-fund/venv/bin/streamlit run dashboard/app.py --server.port=8501 --server.address=0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/var/log/hedge/hedge-dashboard.out.log
stderr_logfile=/var/log/hedge/hedge-dashboard.err.log
environment=PYTHONPATH="/root/hedge-fund",ENV="prod",PYTHONUNBUFFERED="1"
