# deploy/supervisor/hedge-dashboard.conf
[program:hedge-dashboard]
directory=/root/hedge-fund
command=/root/hedge-fund/venv/bin/streamlit run dashboard/multi_strategy_dashboard.py --server.port=8501 --server.address=0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/var/log/hedge-dashboard.log
stderr_logfile=/var/log/hedge-dashboard.err
environment=PYTHONPATH="/root/hedge-fund",ENV="prod",PYTHONUNBUFFERED="1"
